# driver_ipmi/100-udev-driver_measure.rules
KERNEL=="driver_measure", ACTION=="add", MODE="0666", RUN+="/bin/chmod 666 /dev/measure"
